# ec2-cluster/modules/node/templates/l1/l1.service
[Unit]
Description=Node DAG service

[Service]
WorkingDirectory=/home/${user}/tessellation/l1
User=${user}
Group=${user}
ExecStart=/home/${user}/tessellation/l1/start
