# ec2-cluster/modules/node/templates/l0/l0.service
[Unit]
Description=Node service

[Service]
WorkingDirectory=/home/${user}/tessellation/l0
User=${user}
Group=${user}
ExecStart=/home/${user}/tessellation/l0/start
